# crontab
# Run every 15 seconds
*/15 * * * * * * /usr/local/bin/resolv-conf-fixer.sh
